# debian/rustpos-printclient.service
[Unit]
Description=RustPOS Remote Print Client
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/rustpos-printclient
ExecStart=/opt/rustpos-printclient/rustpos-printclient /opt/rustpos-printclient/printclient.toml
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/rustpos-printclient
PrivateTmp=true
SupplementaryGroups=dialout

[Install]
WantedBy=multi-user.target
